# deploy/video-control.service
[Unit]
Description=Control de video GPIO + ffplay (sin entorno gráfico, salida DRM/HDMI)
After=multi-user.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=video1
Group=video1
SupplementaryGroups=video render input
WorkingDirectory=/home/video1/video-player
# Sin X11 (evitar dependencias de entorno gráfico)
Environment=DISPLAY=
# Audio: hdmi (default) o externa — descomentar para placa USB/HAT:
# Environment=AUDIO_SALIDA=externa
# Environment=ALSA_EXTERNA=plughw:1,0
# Environment=ALSA_HDMI=plughw:CARD=vc4hdmi0,DEV=0
# Opcional: configuración de backend ffplay
# Environment=FFPLAY_BIN=ffplay
# Environment=FFPLAY_LOGLEVEL=warning
# Environment=FFPLAY_EXTRA_ARGS=-fflags +genpts
# Esperar a que exista el dispositivo DRM del HDMI
ExecStartPre=/bin/bash -c 'for i in $(seq 1 30); do test -e /dev/dri/card0 && exit 0; sleep 1; done; echo "WARN: /dev/dri/card0 no apareció"; exit 0'
ExecStartPre=/bin/sleep 8
ExecStart=/usr/bin/python3 /home/video1/video-player/video_control.py
Restart=on-failure
RestartSec=8
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
